# eink-calendar.service
[Unit]
Description=E-Ink Calendar Display
After=network.target

[Service]
Type=simple
User=john
WorkingDirectory=/home/john/eink-calendar
ExecStart=/usr/bin/python3 /home/john/eink-calendar/daemon.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
